# deployment/config/supervisord/rag.conf
#================================================================
# RAG知识库系统 - Supervisor 配置
# 适用于: 需要进程管理的环境
#================================================================

[program:rag]
command=java -Xms512m -Xmx1024m -XX:+UseG1GC -jar /opt/rag/backend/rag-1.0.0.jar
directory=/opt/rag/backend
user=rag
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/rag/supervisor.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/rag/supervisor_error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=JAVA_HOME="/usr/lib/jvm/temurin-21-jdk",SPRING_PROFILES_ACTIVE="prod"
stopasgroup=true
killasgroup=true
priority=999
numprocs=1
process_name=%(program_name)s

; 进程监控
;startretries=3
;exitcodes=0
;stopsignal=TERM
;stopwaitsecs=10
